# autostart/vesicolos.desktop
[Desktop Entry]
Type=Application
Name=VESICOLOS
Exec=lxterminal -e '/bin/bash /home/user/Desktop/vesicolos/start.sh'
Terminal=true
