# autostart/vesicolos.desktop
[Desktop Entry]
Type=Application
Name=VESICOLOS
Exec=~/Desktop/vesicolos/start.sh
Terminal=true
